# criminalidad.service.example
[Unit]
Description=Criminalidad España - FastAPI Application
After=network.target postgresql.service

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/criminalidad_app/backend
Environment="PATH=/path/to/criminalidad_app/.venv/bin"
ExecStart=/path/to/criminalidad_app/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
